# finance-bot.service
# =============================================================================
# systemd service — FinanceBot
# =============================================================================
# Скопируйте этот файл на сервер:
#   sudo cp finance-bot.service /etc/systemd/system/finance-bot.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable finance-bot
#   sudo systemctl start finance-bot
#
# Логи:  journalctl -u finance-bot -f
# =============================================================================

[Unit]
Description=FinanceBot Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Finance_bot
EnvironmentFile=/root/Finance_bot/.env
ExecStart=/root/Finance_bot/.venv/bin/python main.py
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
